# installer/build/bootable/root/etc/profile.d/kubernetes.sh
if [ -d /opt/bin ]; then
        pathappend /opt/bin
fi