# installer/build/bootable/root/etc/systemd/system/kubeadm.service
[Unit]
Description=Kubeadm provisioning script
ConditionPathExists=!/etc/vmware/kubeadm.provisioned
After=network-online.service getty@tty2.service

[Service]
Type=oneshot
ExecStart=/opt/ova/bin/kubeadm-provisioning.sh
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/bin

[Install]
WantedBy=multi-user.target
