# installer/build/bootable/root/etc/systemd/system/kubernetes-watchdog.timer
[Unit]
Description=Kubernetes watchdog timer

[Timer]
OnBootSec=0min
OnCalendar=*:0/1
Unit=kubernetes-watchdog.service

[Install]
WantedBy=timers.target