# installer/build/bootable/root/etc/systemd/system/ova-firstboot.service
[Unit]
Description=one time set up script
Before=sshd.service

[Service]
Type=oneshot
ExecStart=/etc/vmware/ova-firstboot.sh

[Install]
WantedBy=network-online.target
